# cam/oscam/debian/oscam.prerm
#!/bin/sh

case "$1" in
    remove)
        systemctl stop oscam > /dev/null 2>&1 || true
        systemctl disable oscam > /dev/null 2>&1 || true
        stop oscam > /dev/null 2>&1 || true
        rm -f /var/log/upstart/oscam* > /dev/null 2>&1 || true
    ;;
esac