# cam/wicard-installer/data/etc/init/wicard.conf
description "Wicard emulator"

#start on runlevel [2345]
#stop  on runlevel  [016]

start on vdr-started
stop  on vdr-stopped

emits wicard-started
emits wicard-stopped

console log

respawn
respawn limit 5 60

script
    exec /usr/bin/wicardd -c /etc/wicard/wicard.conf
end script

post-start script
    initctl emit --no-wait wicard-started
end script

post-stop script
    initctl emit --no-wait wicard-stopped
end script
